# deploy/vibevoice-server.service
[Unit]
Description=VibeVoice ASR Server
After=network-online.target default.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/vibe-voice-vendor
ExecStart=%h/Desktop/vibe-voice-vendor/.venv/bin/python -m server \
    --vllm-base-url http://127.0.0.1:37845 \
    --host 127.0.0.1 \
    --port 54912 \
    --max-audio-bytes 524288000 \
    --max-queue-size 50 \
    --jwt-public-key-file %h/Desktop/vibe-voice-vendor/keys/public.pem \
    --revoked-tokens-file %h/Desktop/vibe-voice-vendor/revoked_tokens.txt \
    --require-https true \
    --vllm-model-name vibevoice \
    --vllm-temperature 0.0 \
    --vllm-top-p 1.0
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
