# deploy/vibevoice-proxy.service
[Unit]
Description=VibeVoice TLS Reverse Proxy
After=network-online.target default.target vibevoice-server.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/Desktop/vibe-voice-vendor/rust_proxy
ExecStart=%h/Desktop/vibe-voice-vendor/rust_proxy/target/release/vvv_proxy \
    --upstream-host 127.0.0.1 \
    --upstream-port 54912 \
    --listen-host 0.0.0.0 \
    --listen-port 42862 \
    --max-body-size 524288000 \
    --cert-path %h/Desktop/vibe-voice-vendor/certs/self-signed/fullchain.pem \
    --key-path %h/Desktop/vibe-voice-vendor/certs/self-signed/privkey.pem \
    --cert-validity-days 3650 \
    --cert-check-interval-secs 3600
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
